# itgmania-bsys-0.1.8-amd64/etc/systemd/system/os-restart-watcher.path
[Unit]
Description=Monitor changes to /var/local/os-restart flagfile

[Path]
PathChanged=/var/local/os-restart

[Install]
WantedBy=multi-user.target
