# itgmania-bsys-0.1.8-amd64/etc/systemd/system/os-shutdown-watcher.path
[Unit]
Description=Monitor changes to /var/local/os-shutdown flagfile

[Path]
PathChanged=/var/local/os-shutdown

[Install]
WantedBy=multi-user.target
