# itgmania-bsys-0.1.8-amd64/etc/systemd/system/update-watcher.path
[Unit]
Description=Monitor changes to /var/local/update flagfile

[Path]
PathChanged=/var/local/update

[Install]
WantedBy=multi-user.target
