# itgmania-bsys-0.1.8-amd64/etc/systemd/user/itgmania-restart-watcher.service
[Unit]
Description=Restart itgmania when the /var/local/itg-restart flagfile changes.

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/libexec/update-itg-restart-flag.sh

[Install]
WantedBy=default.target
