# itgmania-bsys-0.1.8-amd64/etc/systemd/user/itgmania-restart-watcher.path
[Unit]
Description=Monitor changes to /var/local/itg-restart flagfile

[Path]
PathChanged=/var/local/itg-restart

[Install]
WantedBy=default.target
