# itgmania-bsys-0.1.8-amd64/etc/systemd/user/itgmania-reset-watcher.service
[Unit]
Description=Reset itgmania tp default settings when the /var/local/itg-reset flagfile changes.

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/libexec/update-itg-reset-flag.sh

[Install]
WantedBy=default.target
