# itgmania-bsys-0.1.8-amd64/etc/systemd/user/itgmania-stop-watcher.path
[Unit]
Description=Monitor changes to /var/local/itg-stop flagfile

[Path]
PathChanged=/var/local/itg-stop

[Install]
WantedBy=default.target
