# itgmania-bsys-0.1.8-amd64/etc/systemd/user/itgmania-reset-watcher.path
[Unit]
Description=Monitor changes to /var/local/itg-reset flagfile

[Path]
PathChanged=/var/local/itg-reset

[Install]
WantedBy=default.target
