# itgmania-bsys-0.1.8-amd64/etc/systemd/user/obs-streaming-watcher.path
[Unit]
Description=Monitor changes to /var/local/itg-stream-enabled flagfile

[Path]
PathChanged=/var/local/itg-stream-enabled

[Install]
WantedBy=default.target
