# src/etc/systemd/user/itgmania-stop-watcher.path
[Unit]
Description=Monitor changes to /run/user/1000/itg/itg-stop flagfile

[Path]
PathChanged=/run/user/1000/itg/itg-stop

[Install]
WantedBy=default.target
