# src/etc/systemd/user/itgmania.service
[Unit]
Description=ITGMania Service
Requires=graphical-session.target
After=graphical-session.target

[Service]
Type=exec
Restart=on-failure
RestartSec=3s
ExecStart=/usr/bin/bash -c '/opt/itgmania/start.sh'
Nice=-10

[Install]
WantedBy=default.target
