# src/etc/systemd/user/os-restart-watcher.path
[Unit]
Description=Monitor changes to /run/user/1000/itg/os-restart flagfile

[Path]
PathChanged=/run/user/1000/itg/os-restart

[Install]
WantedBy=default.target
