# src/etc/systemd/user/os-shutdown-watcher.path
[Unit]
Description=Monitor changes to /run/user/1000/itg/os-shutdown flagfile

[Path]
PathChanged=/run/user/1000/itg/os-shutdown

[Install]
WantedBy=default.target
