# src/etc/systemd/user/unclutter.timer
[Unit]
Description=Run unclutter every 20 seconds

[Timer]
OnUnitActiveSec=20s
Unit=unclutter.service

[Install]
WantedBy=timers.target