# src/etc/systemd/user/unclutter.service
[Unit]
Description=Hide Mouse Cursor
After=graphical.target

[Service]
ExecStart=/usr/bin/unclutter-xfixes --timeout 0
Restart=no
