# src/etc/systemd/system/os-shutdown-watcher.path
[Unit]
Description=Monitor changes to /run/user/1000/itg/os-shutdown flagfile
After=default.target
After=user-runtime-dir@1000.service
Requires=user-runtime-dir@1000.service

[Path]
PathChanged=/run/user/1000/itg/os-shutdown

[Install]
WantedBy=multi-user.target
