# src/etc/systemd/system/update-watcher.path
[Unit]
Description=Monitor changes to /run/user/1000/itg/update flagfile
After=graphical.target
Requires=graphical.target

[Path]
PathChanged=/run/user/1000/itg/update

[Install]
WantedBy=graphical.target
