# src/etc/systemd/system/update-watcher.path
[Unit]
Description=Monitor changes to /run/itg/update flagfile

[Path]
PathChanged=/run/itg/update

[Install]
WantedBy=multi-user.target
